# salt/roots/salt/helloworld/gunicorn.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid helloworld
setgid helloworld
chdir /home/helloworld

exec .env/bin/gunicorn app:app
